# serverdiscordbot.service
[Unit]
Description=Discord Bot Service
After=network.target

[Service]
Type=simple
User=SERVER_USERNAME
WorkingDirectory=/home/SERVER_USERNAME/discordbot
ExecStart=/usr/bin/node /home/SERVER_USERNAME/discordbot/bot.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
